# apt-mirror/crontab/apt-mirror
# apt-mirror
0 */12 * * *   /usr/bin/apt-mirror



